# dockerfiles/chronograph/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:chronograf]
command=service chronograf start
